# rootfs/etc/services.d/scheduler/run
#!/command/with-contenv bashio
# ==============================================================================
# 
# Home Assistant Add-on: SimpleScheduler
# 
# ==============================================================================

bashio::log.info "Starting service.d [Scheduler]"

exec /usr/bin/scheduler.sh
