# rootfs/etc/services.d/scheduler/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SimpleScheduler
# Runs the scheduler daemon
# ==============================================================================

bashio::log.info "Starting scheduler daemon"

chmod 777 /simplescheduler/*
exec /usr/local/bin/python3 /simplescheduler/scheduler.py
