# rootfs/etc/services.d/scheduler/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SimpleScheduler
# Runs the scheduler daemon
# ==============================================================================

bashio::log.info "Starting scheduler daemon"
cd /simplescheduler
chmod 777 *.py
exec /usr/local/bin/python3 ./scheduler.py
