# rootfs/etc/services.d/scheduler/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SimpleScheduler
# Runs the scheduler daemon
# ==============================================================================

#bashio::log.info "Starting scheduler daemon"

#exec /usr/local/bin/python3 ./scheduler.py
